# etc/xdg/autostart/copy-missing-config-files.desktop
[Desktop Entry]
Type=Application
Name=Copy missing config files
Exec=/bin/sh -c "cp -rnT /etc/skel ~"
Terminal=false
